# chef/site-cookbooks/blog/templates/default/blog-service.conf.erb
description "Blog front end"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

script
  export DW_DB_URL="<%= node['blog']['db']['url'] %>"
  export ETCD_URL="<%= node['blog']['etcd']['url'] %>"
  export HOST_IP=`/usr/bin/iplookup 'eth0 '`

  export GC_OPTS="-XX:+UseParallelOldGC"
  export JAVA_OPTS="-server -Xms256M -Xmx512M"
  export APP_OPTS="-Ddw.database.url=$DW_DB_URL -Ddw.etcd.hosts=$ETCD_URL -Ddw.etcd.publish.hostName=$HOST_IP"
  export APP_JAR="/opt/blog-service/blog-service-<%= node['blog']['service']['version'] %>-shadow.jar"

  /usr/bin/java $JAVA_OPTS $GC_OPTS $APP_OPTS -jar $APP_JAR server /opt/blog-service/prod_config.yml
end script
